# deploy/radstation.service
# /etc/systemd/system/radstation.service
#
# Main RadStation monitoring service.
# Restart behaviour:  always restart, with a 15 s back-off, max 5 times in
# 2 minutes before systemd gives up and the updater's rollback logic takes over.

[Unit]
Description=RadStation Geiger Counter Monitor
Documentation=https://github.com/navengerge-66/GeigerStation
# Wait until the network stack is fully up (Bluetooth serial may need it)
After=network-online.target
Wants=network-online.target
# Ensure the updater does not run at the same time as a restart
After=radstation-updater.service
# If the service crashes repeatedly, stop retrying after 5 attempts in 2 minutes.
# The updater's rollback logic detects the inactive state and takes over.
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=navenger
WorkingDirectory=/home/navenger/radstat

ExecStart=/home/navenger/radstat/venv/bin/python3 -u /home/navenger/radstat/RadStation_v3.py

# Restart unconditionally — covers crashes, OOM kills, and unhandled exceptions.
Restart=always
RestartSec=15

# Capture stdout/stderr to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=radstation

# Credentials — set real values directly on the Pi, never commit them to git
Environment="TELEGRAM_TOKEN=YOUR_TELEGRAM_TOKEN"
Environment="SUPABASE_URL=https://YOUR_PROJECT_REF.supabase.co"
Environment="SUPABASE_KEY=YOUR_SERVICE_ROLE_KEY"

# Niceness: run slightly below default priority so the OS stays responsive
Nice=5

[Install]
WantedBy=multi-user.target
